# GLUT-Schedule.service
[Unit]                     
Description=GLUT Schedule DingTalk Bot                                                                                                  
After=network.target                                          
                                                                                                                                    
[Service]                                                                                                                               
Type=simple
User=root                                                                                                                               
WorkingDirectory=/root/GLUT-Schedule                                             
ExecStart=/root/GLUT-Schedule/.venv/bin/python app.py
Restart=always                         
RestartSec=10                          
                                                                                                                                    
[Install]                                                                                                                               
WantedBy=multi-user.target 